# archoream/mpdris2-rs/mpdris2-rs.service
[Unit]
Description=Music Player Daemon D-Bus Bridge
Wants=mpd.service
After=mpd.service

[Service]
Type=dbus
BusName=org.mpris.MediaPlayer2.mpd
Restart=on-failure
ExecStart=/usr/bin/mpdris2-rs

[Install]
WantedBy=default.target
